# deploy/easypanel/supervisord.conf
[supervisord]
nodaemon=true
logfile=/dev/null
logfile_maxbytes=0
pidfile=/var/lib/seo/run/supervisord.pid
childlogdir=/var/lib/seo/log

[program:api]
directory=/app/backend
command=uvicorn app.main:app --host 0.0.0.0 --port 8000
autostart=true
autorestart=true
startretries=10
stdout_logfile=/dev/fd/1
stdout_logfile_maxbytes=0
stderr_logfile=/dev/fd/2
stderr_logfile_maxbytes=0

[program:worker]
directory=/app/backend
command=celery -A app.workers.celery_app worker --loglevel=INFO --concurrency=2
autostart=true
autorestart=true
startretries=10
stdout_logfile=/dev/fd/1
stdout_logfile_maxbytes=0
stderr_logfile=/dev/fd/2
stderr_logfile_maxbytes=0

[program:beat]
directory=/app/backend
command=celery -A app.workers.celery_app beat --loglevel=INFO --schedule=/var/lib/seo/celery/celerybeat-schedule
numprocs=1
autostart=true
autorestart=true
startretries=10
stdout_logfile=/dev/fd/1
stdout_logfile_maxbytes=0
stderr_logfile=/dev/fd/2
stderr_logfile_maxbytes=0

[program:nginx]
command=/usr/sbin/nginx -g "daemon off;"
autostart=true
autorestart=true
startretries=10
stdout_logfile=/dev/fd/1
stdout_logfile_maxbytes=0
stderr_logfile=/dev/fd/2
stderr_logfile_maxbytes=0
